# systemd/hid.service
[Unit]
Description=Start Hotkey

[Service]
Type=oneshot
ExecStart=/opt/hotkey/hid.sh

[Install]
WantedBy=multi-user.target
